# webling-sync/systemd/webling-sync-members.service
[Unit]
Description=Webling Members Sync to Supabase
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/kraeki/work/jacwohlen.ch/administra/webling-sync
Environment="ENV=prod"
ExecStart=/home/kraeki/work/jacwohlen.ch/administra/webling-sync/run-members.sh
StandardOutput=append:/home/kraeki/work/jacwohlen.ch/administra/webling-sync/logs/members.log
StandardError=append:/home/kraeki/work/jacwohlen.ch/administra/webling-sync/logs/members.log

[Install]
WantedBy=default.target
